# ansible/roles/otto/templates/otto.service.j2
[Unit]
Description=Otto IoT Device Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User={{ otto_user }}
Group={{ otto_group }}
WorkingDirectory={{ otto_home }}
Environment="OTTO_CONFIG={{ otto_config_file }}"
Environment="OTTO_LOG={{ otto_log_file }}"

# ExecStart should point to the actual Otto application binary/script
# Update this based on your Otto deployment method
ExecStart=/usr/local/bin/otto --config {{ otto_config_file }}

# Restart policy
Restart=on-failure
RestartSec=10s

# Logging
StandardOutput=append:{{ otto_log_file }}
StandardError=append:{{ otto_log_file }}

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
